# ks_scripts/snippets/jon_install
# quick 'n dirty JON install script

echo "
#!/bin/bash
# chkconfig: 2345 90 10
# description: installs JON/JOPR/RHQ upon restart of machine

# install postgres 8.2.14 because JON doesn't support postgres < 8.2
rpm -Uvh http://ftp9.us.postgresql.org/pub/mirrors/postgresql/binary/v8.2.14/linux/rpms/redhat/rhel-5-i386/postgresql-8.2.14-1PGDG.rhel5.i386.rpm http://ftp9.us.postgresql.org/pub/mirrors/postgresql/binary/v8.2.14/linux/rpms/redhat/rhel-5-i386/postgresql-libs-8.2.14-1PGDG.rhel5.i386.rpm http://ftp9.us.postgresql.org/pub/mirrors/postgresql/binary/v8.2.14/linux/rpms/redhat/rhel-5-i386/postgresql-server-8.2.14-1PGDG.rhel5.i386.rpm

# get jon and pop it into a new jon user directory
wget $jon_url -O /root/jon.zip
useradd jon -p $6$XugnrnBjnNdMLzNu$6Aj67pB62YBlNCLw/tR9kwmxPHf0/RZu0qPijl77.r7iZhj2p6b1BPleUS.xXJAxCu3rsN3liSXUfcPzAy9dZ1
chown jon /root/jon.zip
mv /root/jon.zip /home/jon

$SNIPPET('jon_postgres_setup')

# unzip jon
su - jon -c 'unzip jon.zip'
su - jon -c 'rm jon.zip'
su - jon -c 'mv jon* jon'

cp /home/jon/jon/bin/rhq-server.sh /etc/init.d

# prepend chkconfig preamble
echo '#!/bin/sh
#chkconfig: 2345 95 20
#description: Jopr Server
#processname: run.sh
RHQ_SERVER_HOME=/home/jon/jon/
RHQ_SERVER_JAVA_HOME=/usr/lib/jvm/jre-1.6.0-openjdk' > /tmp/out 
cat /etc/init.d/rhq-server.sh >> /tmp/out
mv /tmp/out /etc/init.d/rhq-server.sh
chmod 755 /etc/init.d/rhq-server.sh

# rig JON as a service
chkconfig --add rhq-server.sh
chkconfig rhq-server.sh --list
chkconfig --level 5 rhq-server.sh on
chkconfig joninstall off

service rhq-server.sh start
" > /etc/init.d/joninstall
chmod a+x /etc/init.d/joninstall
chkconfig joninstall on 
